# templates/config/backup/database/crontab
# backup datenbases every 3h, verify & purge old backups daily
15 */3 * * * /usr/bin/duply data backup 2>&1
15 4 * * * /usr/bin/duply data backup_verify_purge --force 2>&1
